# core/debian/postinst
#!/bin/sh
set -e

update-rc.d lithiumcore defaults

if which invoke-rc.d >/dev/null 2>&1; then
  invoke-rc.d lithiumcore restart
else
  /etc/init.d/lithiumcore restart
fi

exit 0
